# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/teamproject/bin/uwsgi srv/FC-TeamProject/.config/uwsgi/uwsgi-app.ini